# crontab
CRON_TZ=UTC
1 0 * * * /usr/local/bin/getSeed.sh >> /proc/1/fd/1 2>>/proc/1/fd/2
